# scripts/systemd/sha-agent.service
[Unit]
Description=SHA Endpoint Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/sha-agent -config /etc/sha/agent-config.json -loop
Restart=always
RestartSec=30
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=30
NoNewPrivileges=false
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sha-agent

[Install]
WantedBy=multi-user.target
