# scripts/launchd/com.sha.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sha.agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/sbin/sha-agent</string>
    <string>-config</string>
    <string>/Library/Application Support/SHA/agent-config.json</string>
    <string>-loop</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/var/log/sha-agent.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/sha-agent.err</string>
</dict>
</plist>
